# services/waypoint-scheduler/run
#!/bin/sh
# Waypoint Scheduler - Runit Service
# Calls Rust-based scheduler binary for multi-schedule support

exec 2>&1

# Check if waypoint-scheduler binary is available
if ! command -v waypoint-scheduler >/dev/null 2>&1; then
    echo "ERROR: waypoint-scheduler not found in PATH"
    echo "Please install waypoint before running this service"
    exit 1
fi

# Check if waypoint-cli is available (used by scheduler)
if ! command -v waypoint-cli >/dev/null 2>&1; then
    echo "ERROR: waypoint-cli not found in PATH"
    echo "Please install waypoint before running this service"
    exit 1
fi

# Run the Rust scheduler
exec waypoint-scheduler
